# scripts/scripts/monitor/96-monitor.rules
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/local/bin/hotplug_monitor.sh"

